# s6-rc.d/task-manager-jsapi/run
#!/command/with-contenv sh

# Redis/Valkey connection (shared with orchestrator)
export SUPERBLOCKS_AGENT_REDIS_PASSWORD="${SB_AGENT_REDIS_PASSWORD:-${SUPERBLOCKS_AGENT_REDIS_PASSWORD}}"
export SUPERBLOCKS_WORKER_SANDBOX_TRANSPORT_REDIS_HOST="${SUPERBLOCKS_WORKER_SANDBOX_TRANSPORT_REDIS_HOST:-127.0.0.1}"
export SUPERBLOCKS_WORKER_SANDBOX_TRANSPORT_REDIS_PORT="${SUPERBLOCKS_WORKER_SANDBOX_TRANSPORT_REDIS_PORT:-6379}"
export SUPERBLOCKS_WORKER_SANDBOX_TRANSPORT_REDIS_PASSWORD="${SUPERBLOCKS_WORKER_SANDBOX_TRANSPORT_REDIS_PASSWORD:-${SUPERBLOCKS_AGENT_REDIS_PASSWORD}}"
export SUPERBLOCKS_WORKER_SANDBOX_STORE_REDIS_PASSWORD="${SUPERBLOCKS_WORKER_SANDBOX_STORE_REDIS_PASSWORD:-${SUPERBLOCKS_AGENT_REDIS_PASSWORD}}"

# Sandbox configuration - static mode connecting to local JSAPI sandbox process
export SUPERBLOCKS_WORKER_SANDBOX_SANDBOX_ADDRESS="${SUPERBLOCKS_WORKER_SANDBOX_SANDBOX_ADDRESS_JSAPI:-localhost:50071}"
export SUPERBLOCKS_WORKER_SANDBOX_SANDBOX_LANGUAGE="javascript"

# Worker configuration
export SUPERBLOCKS_WORKER_SANDBOX_WORKER_GROUP="${SUPERBLOCKS_WORKER_SANDBOX_WORKER_GROUP:-main}"
export SUPERBLOCKS_WORKER_SANDBOX_WORKER_CONSUMER_GROUP="${SUPERBLOCKS_WORKER_SANDBOX_WORKER_CONSUMER_GROUP_JSAPI:-sandbox-javascriptsdkapi-workers}"
export SUPERBLOCKS_WORKER_SANDBOX_WORKER_BUCKETS="${SUPERBLOCKS_WORKER_SANDBOX_WORKER_BUCKETS:-BA}"
export SUPERBLOCKS_WORKER_SANDBOX_WORKER_PLUGINS="${SUPERBLOCKS_WORKER_SANDBOX_WORKER_PLUGINS_JSAPI:-javascriptsdkapi}"
export SUPERBLOCKS_WORKER_SANDBOX_WORKER_EVENTS="${SUPERBLOCKS_WORKER_SANDBOX_WORKER_EVENTS:-execute metadata test pre_delete}"

# gRPC server ports (different from JS to avoid conflicts)
export SUPERBLOCKS_WORKER_SANDBOX_VARIABLE_STORE_GRPC_PORT="${SUPERBLOCKS_WORKER_SANDBOX_VARIABLE_STORE_GRPC_PORT_JSAPI:-50070}"
export SUPERBLOCKS_WORKER_SANDBOX_VARIABLE_STORE_HTTP_PORT="${SUPERBLOCKS_WORKER_SANDBOX_VARIABLE_STORE_HTTP_PORT_JSAPI:-8087}"
export SUPERBLOCKS_WORKER_SANDBOX_STREAMING_PROXY_GRPC_PORT="${SUPERBLOCKS_WORKER_SANDBOX_STREAMING_PROXY_GRPC_PORT_JSAPI:-50072}"

# Observability
export SUPERBLOCKS_WORKER_SANDBOX_LOG_LEVEL="${SUPERBLOCKS_WORKER_SANDBOX_LOG_LEVEL:-${SUPERBLOCKS_AGENT_LOG_LEVEL:-info}}"
export SUPERBLOCKS_WORKER_SANDBOX_OTEL_COLLECTOR_HTTP_URL="${SUPERBLOCKS_WORKER_SANDBOX_OTEL_COLLECTOR_HTTP_URL:-http://127.0.0.1:4318}"

# In the OPA, sandbox and task-manager are co-located; no Kubernetes pod IP registration
export SUPERBLOCKS_WORKER_SANDBOX_IP_FILTER_DISABLED="${SUPERBLOCKS_WORKER_SANDBOX_IP_FILTER_DISABLED:-true}"

# Agent key for authentication
export SUPERBLOCKS_WORKER_SANDBOX_SUPERBLOCKS_KEY="${SUPERBLOCKS_AGENT_KEY:-${SB_AGENT_KEY:-dev-agent-key}}"

echo "{\"level\":\"info\",\"ts\":$(date +%s%3N),\"msg\":\"starting component task-manager-jsapi\",\"component\":\"s6\"}"
s6-setuidgid 1000:1000 /app/task-manager/bin
