# s6-rc.d/workerpy/run
#!/command/with-contenv sh

export SUPERBLOCKS_AGENT_KEY="${SB_AGENT_KEY:-${SUPERBLOCKS_AGENT_KEY}}"
export SUPERBLOCKS_WORKER_EXECUTION_ENV_INCLUSION_LIST="${SB_EXECUTION_ENV_INCLUSION_LIST:-${SUPERBLOCKS_EXECUTION_ENV_INCLUSION_LIST}}"
export SB_METRICS_PORT=19090

export SUPERBLOCKS_PYTHON_EXECUTION_BUILTINS_DENY_LIST="${SB_PYTHON_EXECUTION_BUILTINS_DENY_LIST:-${SUPERBLOCKS_PYTHON_EXECUTION_BUILTINS_DENY_LIST}}"
export SUPERBLOCKS_PYTHON_EXECUTION_IMPORT_DENY_LIST="${SB_PYTHON_EXECUTION_IMPORT_DENY_LIST:-${SUPERBLOCKS_PYTHON_EXECUTION_IMPORT_DENY_LIST}}"

# Set SUPERBLOCKS_PYTHON_EXECUTION_IMPORT_ALLOW_LIST from SB_PYTHON_EXECUTION_IMPORT_ALLOW_LIST or SUPERBLOCKS_PYTHON_EXECUTION_IMPORT_ALLOW_LIST
# If neither are set, we want to ensure SUPERBLOCKS_PYTHON_EXECUTION_IMPORT_ALLOW_LIST remains unset since a value of "" has meaning (reject all imports)
if [ -n "${SB_PYTHON_EXECUTION_IMPORT_ALLOW_LIST+x}" ]; then
  export SUPERBLOCKS_PYTHON_EXECUTION_IMPORT_ALLOW_LIST="$SB_PYTHON_EXECUTION_IMPORT_ALLOW_LIST"
elif [ -n "${SUPERBLOCKS_PYTHON_EXECUTION_IMPORT_ALLOW_LIST+x}" ]; then
  export SUPERBLOCKS_PYTHON_EXECUTION_IMPORT_ALLOW_LIST="$SUPERBLOCKS_PYTHON_EXECUTION_IMPORT_ALLOW_LIST"
fi


export SUPERBLOCKS_AGENT_REDIS_PASSWORD="${SB_AGENT_REDIS_PASSWORD:-${SUPERBLOCKS_AGENT_REDIS_PASSWORD}}"
export SUPERBLOCKS_AGENT_REDIS_KVSTORE_TOKEN="${SUPERBLOCKS_AGENT_REDIS_KVSTORE_TOKEN:-${SUPERBLOCKS_AGENT_REDIS_PASSWORD}}"
export SUPERBLOCKS_AGENT_REDIS_TOKEN="${SUPERBLOCKS_AGENT_REDIS_TOKEN:-${SUPERBLOCKS_AGENT_REDIS_PASSWORD}}"

echo "{\"level\":\"info\",\"ts\":$(date +%s%3N),\"msg\":\"starting component worker.py\",\"component\":\"s6\"}"
python3 /app/worker.py/src/entry.py
