# s6-rc.d/orchestrator/run
#!/command/with-contenv sh

export SUPERBLOCKS_ORCHESTRATOR_SUPERBLOCKS_KEY="${SB_AGENT_KEY:-${SUPERBLOCKS_ORCHESTRATOR_SUPERBLOCKS_KEY}}"
export SUPERBLOCKS_ORCHESTRATOR_AGENT_HOST_URL="${SB_AGENT_URL:-${SUPERBLOCKS_ORCHESTRATOR_AGENT_HOST_URL}}"
export SUPERBLOCKS_ORCHESTRATOR_LOG_LEVEL="${SUPERBLOCKS_ORCHESTRATOR_LOG_LEVEL:-${SUPERBLOCKS_AGENT_LOG_LEVEL}}"
export SUPERBLOCKS_ORCHESTRATOR_AUTH_JWT_ENABLED="${SB_AUTH_JWT_ENABLED:-${SUPERBLOCKS_ORCHESTRATOR_AUTH_JWT_ENABLED}}"
export SUPERBLOCKS_ORCHESTRATOR_SUPERBLOCKS_URL="${SB_URL:-${SUPERBLOCKS_ORCHESTRATOR_SUPERBLOCKS_URL}}"
export SUPERBLOCKS_ORCHESTRATOR_ZEN="${SB_ZEN}"

export SUPERBLOCKS_AGENT_REDIS_PASSWORD="${SB_AGENT_REDIS_PASSWORD:-${SUPERBLOCKS_AGENT_REDIS_PASSWORD}}"
export SUPERBLOCKS_ORCHESTRATOR_STORE_REDIS_PASSWORD="${SUPERBLOCKS_ORCHESTRATOR_STORE_REDIS_PASSWORD:-${SUPERBLOCKS_AGENT_REDIS_PASSWORD}}"
export SUPERBLOCKS_ORCHESTRATOR_TRANSPORT_REDIS_PASSWORD="${SUPERBLOCKS_ORCHESTRATOR_TRANSPORT_REDIS_PASSWORD:-${SUPERBLOCKS_AGENT_REDIS_PASSWORD}}"

echo "{\"level\":\"info\",\"ts\":$(date +%s%3N),\"msg\":\"starting component orchestrator\",\"component\":\"s6\"}"
s6-setuidgid 1000:1000 /app/orchestrator/bin
