# s6-rc.d/sandbox-jsapi/run
#!/command/with-contenv sh

# gRPC server port (where task-manager-jsapi connects to this sandbox)
export SUPERBLOCKS_WORKER_SANDBOX_TRANSPORT_GRPC_PORT="${SUPERBLOCKS_WORKER_SANDBOX_TRANSPORT_GRPC_PORT_JSAPI:-50071}"

# Variable store and streaming proxy addresses (task-manager-jsapi endpoints)
export SUPERBLOCKS_WORKER_SANDBOX_TRANSPORT_VARIABLE_STORE_ADDRESS="${SUPERBLOCKS_WORKER_SANDBOX_TRANSPORT_VARIABLE_STORE_ADDRESS_JSAPI:-localhost:50070}"
export SUPERBLOCKS_WORKER_SANDBOX_TRANSPORT_VARIABLE_STORE_HTTP_ADDRESS="${SUPERBLOCKS_WORKER_SANDBOX_TRANSPORT_VARIABLE_STORE_HTTP_ADDRESS_JSAPI:-localhost:8087}"
export SUPERBLOCKS_WORKER_SANDBOX_TRANSPORT_STREAMING_PROXY_ADDRESS="${SUPERBLOCKS_WORKER_SANDBOX_TRANSPORT_STREAMING_PROXY_ADDRESS_JSAPI:-localhost:50072}"

# Logging
export SUPERBLOCKS_WORKER_SANDBOX_LOG_LEVEL="${SUPERBLOCKS_WORKER_SANDBOX_LOG_LEVEL:-${SUPERBLOCKS_AGENT_LOG_LEVEL:-info}}"

# Required env vars (not used at runtime but expected at import time by worker.js)
export SUPERBLOCKS_AGENT_KEY="${SUPERBLOCKS_AGENT_KEY:-}"
export SUPERBLOCKS_TUNNEL_PRIVATE_KEY_RSA="${SUPERBLOCKS_TUNNEL_PRIVATE_KEY_RSA:-}"
export SUPERBLOCKS_TUNNEL_PRIVATE_KEY_ED25519="${SUPERBLOCKS_TUNNEL_PRIVATE_KEY_ED25519:-}"

echo "{\"level\":\"info\",\"ts\":$(date +%s%3N),\"msg\":\"starting component sandbox-jsapi\",\"component\":\"s6\"}"
node /app/sandbox-js/dist/index.js
